# testing/scylla/etc/supervisord.conf.d/scylla-server.conf
[program:scylla]
command=/scylla-service.sh
stdout_logfile=/var/log/scylla/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/scylla/stderr
stderr_logfile_maxbytes=0