# docker/scylla/scylla-manager-agent.conf
[program:scylla-manager-agent]
command=/scylla-manager-agent-docker-entrypoint.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0