# dist/systemd/scylla-manager-agent.service
[Unit]
Description=Scylla Manager Agent
Wants=scylla-server.service
After=scylla-server.service
After=network-online.target

[Service]
Type=simple
User=scylla-manager
Group=scylla-manager
EnvironmentFile=/etc/scylla-manager-agent/credentials.env
ExecStart=/usr/bin/scylla-manager-agent
TimeoutStartSec=900
KillMode=process
Restart=on-failure
Slice=scylla-helper

[Install]
WantedBy=multi-user.target
