# dist/systemd/scylla-manager-agent.service
[Unit]
Description=Scylla Manager Agent
Wants=scylla-server.service
Wants=scylla-helper.slice
After=scylla-server.service
After=network-online.target

[Service]
Type=simple
User=scylla-manager
Group=scylla-manager
ExecStart=/usr/bin/scylla-manager-agent
TimeoutStartSec=900
KillMode=process
Restart=on-failure
RestartPreventExitStatus=1
StandardOutput=journal
StandardError=journal
SyslogLevelPrefix=false
Slice=scylla-helper.slice

[Install]
WantedBy=multi-user.target
