# dist/systemd/scylla-mgmt.service
[Unit]
Description=Scylla Management Server
After=network-online.target
Wants=scylla-server.service

[Service]
Type=simple
User=scylla
Group=scylla
ExecStart=/usr/bin/scylla-mgmt server
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
