# dist/systemd/scylla-mgmt.service
[Unit]
Description=Scylla Management Server
Wants=scylla-server.service
After=scylla-server.service
After=network-online.target

[Service]
Type=simple
User=scylla-mgmt
Group=scylla-mgmt
ExecStart=/usr/bin/scylla-mgmt
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
