# dist/systemd/scylla-manager.service
[Unit]
Description=Scylla Manager Server
Wants=scylla-server.service
After=scylla-server.service
After=network-online.target

[Service]
Type=simple
User=scylla-manager
Group=scylla-manager
ExecStart=/usr/bin/scylla-manager
TimeoutStartSec=900
KillMode=process
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
