# dist/systemd/scylla-manager.service
[Unit]
Description=Scylla Manager Server
Wants=scylla-server.service
After=scylla-server.service
After=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=scylla-manager
Group=scylla-manager
ExecStart=/usr/bin/scylla-manager
LimitNOFILE=100000
TimeoutStartSec=900
KillMode=process
Restart=on-failure
StandardOutput=journal
StandardError=journal
SyslogLevelPrefix=false

[Install]
WantedBy=multi-user.target
